# docker/geth/service/ethnetintel/run
#!/bin/sh
# Copyright 2017 - Dechain User Group
exec 2>&1

exec su -c "/home/eth/gethStartIntel.sh"

